# backend/freemarket.service
# Gunicorn configuration for FreeMarket
# Save as /etc/systemd/system/freemarket.service

[Unit]
Description=FreeMarket FastAPI application
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/freemarket/backend
Environment=PATH=/opt/freemarket/backend/venv/bin
ExecStart=/opt/freemarket/backend/venv/bin/gunicorn main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/freemarket/access.log \
    --error-logfile /var/log/freemarket/error.log
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freemarket

[Install]
WantedBy=multi-user.target
