# backend/freemarket-bot.service
# Telegram Bot service for FreeMarket
# Save as /etc/systemd/system/freemarket-bot.service

[Unit]
Description=FreeMarket Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/freemarket/backend
Environment=PATH=/opt/freemarket/backend/venv/bin
ExecStart=/opt/freemarket/backend/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freemarket-bot

[Install]
WantedBy=multi-user.target
